# wakeup/etc/init/network-wakeup.conf
description "wakeup other computers"
author "Marco Scholl <yavdr@marco-schollde>"

start on net-device-up IFACE=*[0123456789]

task

script
  exec yavdr-network-wakeup
end script
