# base/etc/init/setup-wol.conf
description "set wol flag"
author "Marco Scholl <yavdr@marco-scholl.de>"

start on net-device-up IFACE=*[0123456789]

task

env IFACE

script
  exec ethtool -s $IFACE wol g
end script
